# salt/mac2hostname/upstart.conf
# Upstart Configuration

description     "mac to hostname web app"
author          "Flavio Castelli <flavio@castelli.me>"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])

stop on runlevel [016]

respawn

exec /usr/bin/mac2hostname
